# caddy/caddy.service.d/override.conf
[Service]
EnvironmentFile=&HOME&/caddy/conf.env